# etc/profile.d/android.bash
ANDROID_SDK_PATH=~/opt/android
add_path $ANDROID_SDK_PATH/tools
add_path $ANDROID_SDK_PATH/platform-tools
