# etc/profile.d/autojump.bash
../../Code/forked/autojump/autojump.bash