# etc/profile.d/go.bash
# Google Go compilation env variables
export GOROOT=$HOME/Code/forked/go
export GOBIN=$GOROOT/bin
add_path $GOBIN
