# etc/profile.d/ledger.bash
if [ -f "$HOME/Documents/Regnskab/ledger.dat" ]; then
    export LEDGER="$HOME/Documents/Regnskab/ledger.dat"
fi
