# etc/profile.d/rust.bash
add_path ~/opt/rust-nightly/bin
add_path ~/opt/rust-nightly/lib LD_LIBRARY_PATH
