# etc/profile.d/texlive.bash
# TeXLive setup
add_path /usr/local/texlive/2010/bin/x86_64-linux
alias tlmgr='sudo /usr/local/texlive/2010/bin/x86_64-linux/tlmgr'
