# etc/profile.d/texlive.bash
# TeXLive setup
add_path ~/opt/texlive/2012/bin/x86_64-linux
