# etc/profile.d/z.bash
_Z_DATA="$HOME/.local/share/z/z.db"

. ~/Code/cloned/z/z.sh
